# cron/cron.d/mydumper
SHELL=/bin/sh

7 */1 * * * root /run-mydumper.sh
